# system_files/etc/systemd/system/org.45drives.Scheduler.service
[Unit]
Description=45Drives Scheduler D-Bus Service
After=network.target dbus.service
Requires=dbus.service

[Service]
Type=dbus
BusName=org.houston.Scheduler

# Run Python unbuffered so logs flush immediately and systemd tracks the real PID
ExecStart=/usr/bin/python3 -u /usr/local/libexec/45drives-schedulerd

User=root
Group=root

# Faster, cleaner restarts
KillSignal=SIGINT
TimeoutStopSec=10s
Restart=on-failure
RestartSec=250ms
KillMode=control-group

# Hardening (kept), but allow the paths this service genuinely writes to.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
# We need to write under /var/lib/houston (state), /run/user (user manager runtime links),
# and (for legacy/system tasks) /etc/systemd/system.
ReadWritePaths=/var/lib/houston /run/user /etc/systemd/system

[Install]
WantedBy=multi-user.target
