# packaging/ubuntu-jammy/postinst
#!/bin/bash
set -e
export DEBIAN_FRONTEND=noninteractive

# Toggle daemon enablement for this build
# 0 = legacy-only (no daemon), 1 = manage and start daemon if installed
ENABLE_DAEMON=0

# Create state dir (idempotent)
install -d -m 0755 /var/lib/45drives/houston/scheduler

if [ "$ENABLE_DAEMON" = "1" ]; then
    # Reload systemd; enable & start the D-Bus service
    systemctl daemon-reexec || true
    systemctl daemon-reload || true

    if systemctl cat org.45drives.Scheduler.service &>/dev/null; then
        systemctl enable --now org.45drives.Scheduler.service || true
    fi
fi

exit 0
