# system_files/opt/45drives/houston/scheduler/templates/Schedule.timer
[Unit]
Description={description}

[Timer]
{on_calendar_lines}
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target