# system_files/opt/45drives/houston/scheduler/templates/Task.service
[Unit]
Description=Service for {task_name}
StartLimitBurst=2
StartLimitInterval=80s

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile={env_path}
ExecStart={ExecStart}
Restart=on-failure
RestartSec=30sec

[Install]
WantedBy=multi-user.target